# deploy/systemd/ai-control-plane-cert-renewal.service.tmpl
# AI Control Plane - Scheduled Certificate Renewal Service Template
# =============================================================================
# Purpose: Run one threshold-based certificate renewal evaluation.
# Responsibilities:
#   - Execute the typed certificate renewal workflow.
#   - Preserve renewal rollback artifacts when ACP forces reissuance.
# Scope:
#   - Host-first certificate lifecycle automation only.
# Usage:
#   - Render this template and install it to /etc/systemd/system/<service>-cert-renewal.service
# Invariants/Assumptions:
#   - Uses the canonical typed CLI path for certificate operations.
# =============================================================================

[Unit]
Description=AI Control Plane certificate renewal job
Documentation=https://github.com/fitchmultz/ai-control-plane/blob/main/docs/deployment/CERTIFICATE_LIFECYCLE.md
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
User={{SERVICE_USER}}
Group={{SERVICE_GROUP}}
WorkingDirectory={{WORKING_DIR}}
EnvironmentFile={{ENV_FILE}}
Environment=ACP_REPO_ROOT={{WORKING_DIR}}
ExecStart=/bin/bash -lc 'set -euo pipefail; ./scripts/acpctl.sh cert renew --threshold-days {{RENEW_THRESHOLD_DAYS}}'
TimeoutStartSec=600
NoNewPrivileges=yes
UMask=0077
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{SERVICE_NAME}}-cert-renewal
